# generic/huawei-wmi-1.2.0/huawei-wmi-privilege.service
[Unit]
Description=Huawei WMI Privilege
After=sys-devices-platform-huawei\x2dwmi.device

[Service]
Type=oneshot
ExecStart=/usr/lib/huawei-wmi/privilege
