# generic/huawei-wmi-1.1.1/huawei-wmi-privilege.service
[Unit]
Description=Huawei WMI Privilege
After=sys-devices-platform-huawei\x2dwmi.device

[Service]
Type=oneshot
ExecStart=/usr/lib/huawei-wmi/privilege

[Install]
WantedBy=sys-devices-platform-huawei\x2dwmi.device
